# FusionIDS-main/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[program:target]
command=python3 /app/target_server.py
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/var/log/target.log
stderr_logfile=/var/log/target.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=1

[program:blockchain]
command=python3 run_node.py
directory=/app/updated_blockchain
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/blockchain.log
stderr_logfile=/var/log/blockchain.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=2

[program:mlserver]
command=python3 -m uvicorn inference.ml_server:app --host 0.0.0.0 --port 8000
directory=/app/updated_model
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/mlserver.log
stderr_logfile=/var/log/mlserver.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=3

; Exclude internal ports:
;   port 5000 — blockchain PBFT traffic between nodes
;   port 8000 — ML server internal calls
; Only capture real external attack traffic
[program:capturer]
command=/app/traffic_capturer_updated/build/ids eth0 --filter "tcp and not port 5000 and not port 8000"
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/var/log/capturer.log
stderr_logfile=/var/log/capturer.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
priority=4
